# templates/default/nginx.service.erb
[Unit]
Description=The nginx HTTP and reverse proxy server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
<% if node['openresty']['LimitNOFILE'] %>
LimitNOFILE=<%= node['openresty']['LimitNOFILE'] %>
<% end -%>
PIDFile=<%= @pid %>
ExecStartPre=<%= @src_binary %> -t -c <%= node['openresty']['source']['conf_path'] %>
ExecStart=<%= @src_binary %> -c <%= node['openresty']['source']['conf_path'] %>
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=process
KillSignal=SIGQUIT
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
